# evil-assistant.service
[Unit]
Description=Evil Assistant - Demonic Voice Assistant
Documentation=https://github.com/jkinman/evilassistant
After=network-online.target sound.target pulseaudio.service
Wants=network-online.target
RequiresMountsFor=/home

[Service]
Type=simple
User=jkinman
Group=jkinman
SupplementaryGroups=audio
WorkingDirectory=/home/jkinman/evilassistant

# Environment setup
Environment=PATH=/home/jkinman/evilassistant/.venv/bin
Environment=PYTHONPATH=/home/jkinman/evilassistant
Environment=PHILIPS_HUE_BRIDGE_IP=192.168.0.20
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Main command
ExecStart=/home/jkinman/evilassistant/.venv/bin/python -m evilassistant --transcription

# Process management
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security and resource limits
NoNewPrivileges=true
MemoryMax=2G
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=evil-assistant

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
